# files/system/usr/lib/systemd/system/uupd.timer.d/weekly.conf
[Timer]
# Override default daily schedule — run weekly on Thursday after Aurora's stable release
OnCalendar=
OnCalendar=Thu 04:00:00
